# rootfs/etc/s6-overlay/s6-rc.d/cameraui/run
#!/command/with-contenv bash
echo "[cameraui] starting (home: ${CAMERA_UI_HOME_PATH})"
cd /root
ulimit -c 0
exec cameraui run --allow-root --docker -H "${CAMERA_UI_HOME_PATH}"
