# rootfs/etc/s6-overlay/s6-rc.d/dbus/run
#!/command/with-contenv bash
if [[ "${CAMERAUI_DOCKER_AVAHI}" != "true" ]]; then
    echo "[dbus] CAMERAUI_DOCKER_AVAHI != true — dbus disabled"
    exec sleep infinity
fi
echo "[dbus] starting system bus"
exec dbus-daemon --system --nofork --nopidfile
